# thirdparty/stud/stud.service
[Unit]
Description=stud - The Scalable TLS Unwrapping Daemon
Documentation=https://github.com/bumptech/stud
After=network.target
Wants=network.target

[Service]
Type=simple
User=stud
Group=stud
NoNewPrivileges=yes

# stud binary location
ExecStart=/usr/local/bin/stud /etc/stud/stud.pem

# Basic configuration (override with ExecStart options)
# Frontend (listening) address - default: *:8443
# Backend (upstream) address - default: 127.0.0.1:8000
# TLS version enforcement - uncomment and modify as needed:
# --tls11 (minimum TLS 1.1)
# --tls12 (minimum TLS 1.2)
# --tls13 (minimum TLS 1.3)

# Example with TLS 1.2 enforcement:
# ExecStart=/usr/local/bin/stud --tls12 --frontend *,8443 --backend 127.0.0.1,80 /etc/stud/stud.pem

# Additional security options
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/var/run/stud /var/log
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
MemoryDenyWriteExecute=yes

# Resource limits
LimitNOFILE=65536

# Restart on failure
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stud

[Install]
WantedBy=multi-user.target